# docker/all-in-one/rootfs/etc/s6-overlay/s6-rc.d/worker/run
#!/bin/bash
# BullMQ-Worker: Transcode + Import + Caption + Moderation + Scheduler.
set -e
cd /app/apps/worker
exec pnpm start
